# crontab
0 9 * * * python /usr/src/reports_data_import/generate_employee_report.py >> /var/log/cron/cron.log 2>&1
